# resources/systemd/tms-prod/ewms-tms-prod-watcher.path
[Path]
PathModified=/home/ewms/resources/systemd/tms-prod/apptainer_container_symlink

[Install]
WantedBy=multi-user.target
