# resources/systemd/tms-dev/ewms-tms-dev-watcher.path
[Path]
PathModified=/home/ewms/resources/systemd/tms-dev/envfile

[Install]
WantedBy=multi-user.target
