# resources/systemd/tms-dev/ewms-tms-dev-watcher.path
[Path]
PathModified=/home/ewms/tms-dev/apptainer_container_symlink

[Install]
WantedBy=multi-user.target
